# iApache2.service.template
[Unit]
Description=APACHE2 GUI             

[Service]
Type=simple
ExecStart=/bin/bash <P_PATH>/.run.sh

[Install]
WantedBy=multi-user.target