# iApache2.service.template
[Unit]
Description=APACHE2 GUI             

[Service]
Type=simple
ExecStart=cd <P_PATH> && npm start

[Install]
WantedBy=multi-user.target